# homeassistant/old/scripts/ruuvitag-manager/ruuvitag-manager.service
[Unit]
Description=RuuviTag Manager - reads environmental data from RuuviTag BLE sensors and publishes to MQTT
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/homeassistant/.homeassistant/scripts/ruuvitag-manager/manager.py
Restart=on-failure
RestartSec=5
User=homeassistant
Group=homeassistant

# Environment variables
Environment=PYTHONUNBUFFERED=1
Environment=RUUVI_BLE_ADAPTER=bleak

# Capabilities needed for Bluetooth access
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ruuvitag-manager

[Install]
WantedBy=multi-user.target
